# sensors-to-mqtt.service
[Unit]
Description=Sensors to MQTT Bridge Service
Documentation=https://github.com/askrejans/sensors-to-mqtt
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=sensors
Group=sensors

# Path to the binary
ExecStart=/usr/bin/sensors-to-mqtt --config /etc/sensors-to-mqtt/settings.toml

# Working directory (app writes no state; /tmp is always available)
WorkingDirectory=/tmp

# Restart policy
Restart=always
RestartSec=10

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

# Hardware device access for the service user.
# i2c    covers /dev/i2c-* (I2C sensors: MPU-6500, BMP280, BME280, SHT31, BH1750, INA219, ADS1115)
# dialout covers /dev/ttyACM*, /dev/ttyUSB*, /dev/ttyS* (serial sensors: SDS011, etc.)
# gpio   covers /dev/gpiomem, /dev/gpiochip* (GPIO sensors)
SupplementaryGroups=i2c dialout gpio

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallErrorNumber=EPERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sensors-to-mqtt

# Environment variables (optional - can also use settings.toml)
# Environment="SENSORS_TO_MQTT__LOG_LEVEL=info"
# Environment="SENSORS_TO_MQTT__MQTT__HOST=mqtt.local"

[Install]
WantedBy=multi-user.target
